# fan_control.service
[Unit]
Description=Fan Control
After=network.target

[Service]
ExecStart=/opt/fan_control
Restart=always
